# infra/ansible/templates/idcd-agent.service.j2
# /etc/systemd/system/idcd-agent.service
# Generated by Ansible — do NOT edit manually.

[Unit]
Description=idcd Agent
Documentation=https://github.com/kite365/idcd
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ agent_user }}
Group={{ agent_group }}

ExecStart={{ agent_bin_path }} --config {{ agent_config_path }}
Restart=on-failure
RestartSec=10
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Environment
Environment=AGENT_DATA_DIR={{ agent_data_dir }}
Environment=IDCD_CONFIG={{ agent_config_path }}

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ agent_data_dir }} /var/log/idcd-agent
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# CAP_NET_RAW required for ICMP ping probes
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW

# Resource limits
LimitNOFILE=65536
MemoryMax=512M
TasksMax=100

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=idcd-agent

[Install]
WantedBy=multi-user.target
